# dist/radix-gateway/etc/profile.d/radix-gateway.sh
FOPIS_PATH=/opt/fopis
export FOPIS_PATH="$FOPIS_PATH"
export PATH="$PATH:/opt/fopis/bin"
